# roles/apt/templates/systemd/apt-daily.timer.d/override.conf.j2
{{ ansible_managed | comment }}

[Timer]
OnCalendar=
OnCalendar={{ apt_daily_timer_on_calendar }}
RandomizedDelaySec={{ apt_daily_timer_randomized_delay_sec }}
